# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx

[program:php]
command=php-fpm7.0